# recipes-core/customization-ladybird/files/weston.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/weston --config /etc/xdg/weston/weston-ladybird.ini


